# .config/systemd/user/cortile.service
[Unit]
Description=window tiling manager
After=graphical.target

[Service]
ExecStart=/usr/local/bin/cortile

[Install]
WantedBy=default.target
