# .config/systemd/user/gopls.service
[Unit]
Description=go language server.
After=graphical.target

[Service]
ExecStart=/usr/bin/gopls serve

[Install]
WantedBy=default.target
