# .config/systemd/user/xbanish.service
[Unit]
Description=xbanish

[Service]
ExecStart=/usr/local/bin/xbanish

[Install]
WantedBy=multi-user.target
